# deployment/systemd/option-chain-dashboard.service
# Systemd service for Streamlit dashboard
# For Oracle Cloud / VPS deployment
# Location: /etc/systemd/system/option-chain-dashboard.service

[Unit]
Description=Option Chain Streamlit Dashboard
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/option-chain/OptionChainUsingUpstock
Environment="PATH=/home/ubuntu/option-chain/OptionChainUsingUpstock/venv/bin"
Environment="TZ=Asia/Kolkata"
ExecStart=/home/ubuntu/option-chain/OptionChainUsingUpstock/venv/bin/streamlit run optionchain.py --server.port 8502 --server.headless true

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/ubuntu/option-chain/OptionChainUsingUpstock/logs/dashboard.log
StandardError=append:/home/ubuntu/option-chain/OptionChainUsingUpstock/logs/dashboard-error.log

[Install]
WantedBy=multi-user.target
